# master/ml_screen/var/lib/dpkg/info/raspi-config.list
/.
/etc
/etc/default
/etc/default/cpu_governor
/etc/init.d
/etc/init.d/raspi-config
/etc/sudoers.d
/etc/sudoers.d/010_proxy
/usr
/usr/bin
/usr/bin/raspi-config
/usr/lib
/usr/lib/raspi-config
/usr/lib/raspi-config/cmstart.sh
/usr/lib/raspi-config/init_resize.sh
/usr/share
/usr/share/doc
/usr/share/doc/raspi-config
/usr/share/doc/raspi-config/changelog.gz
/usr/share/doc/raspi-config/copyright
/usr/share/raspi-config
/usr/share/raspi-config/10-blanking.conf
