# master/ml_screen/var/lib/dpkg/info/rpcbind.list
/.
/etc
/etc/default
/etc/default/rpcbind
/etc/init.d
/etc/init.d/rpcbind
/etc/insserv.conf.d
/etc/insserv.conf.d/rpcbind
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/rpcbind.service
/lib/systemd/system/rpcbind.socket
/sbin
/sbin/rpcbind
/usr
/usr/bin
/usr/bin/rpcinfo
/usr/lib
/usr/lib/tmpfiles.d
/usr/lib/tmpfiles.d/rpcbind.conf
/usr/sbin
/usr/share
/usr/share/doc
/usr/share/doc/rpcbind
/usr/share/doc/rpcbind/NEWS.Debian.gz
/usr/share/doc/rpcbind/README.Debian
/usr/share/doc/rpcbind/changelog.Debian.armhf.gz
/usr/share/doc/rpcbind/changelog.Debian.gz
/usr/share/doc/rpcbind/changelog.gz
/usr/share/doc/rpcbind/copyright
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/rpcbind.8.gz
/usr/share/man/man8/rpcinfo.8.gz
/lib/systemd/system/portmap.service
/usr/sbin/rpcinfo
