# master/ml_screen/var/lib/dpkg/info/userconf-pi.list
/.
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/userconfig.service
/usr
/usr/bin
/usr/bin/cancel-rename
/usr/bin/rename-user
/usr/lib
/usr/lib/userconf-pi
/usr/lib/userconf-pi/userconf
/usr/lib/userconf-pi/userconf-service
/usr/share
/usr/share/doc
/usr/share/doc/userconf-pi
/usr/share/doc/userconf-pi/changelog.gz
/usr/share/doc/userconf-pi/copyright
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/userconf-pi
/usr/share/userconf-pi
/usr/share/userconf-pi/sshd_banner
